# contrib/init/menabitcoind.conf
description "MENA Bitcoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env MENABITCOIND_BIN="/usr/bin/menabitcoind"
env MENABITCOIND_USER="menabitcoin"
env MENABITCOIND_GROUP="menabitcoin"
env MENABITCOIND_PIDDIR="/var/run/menabitcoind"
# upstart can't handle variables constructed with other variables
env MENABITCOIND_PIDFILE="/var/run/menabitcoind/menabitcoind.pid"
env MENABITCOIND_CONFIGFILE="/etc/menabitcoin/menabitcoin.conf"
env MENABITCOIND_DATADIR="/var/lib/menabitcoind"

expect fork

respawn
respawn limit 5 120
kill timeout 600

pre-start script
    # this will catch non-existent config files
    # menabitcoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$MENABITCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run menabitcoind."
        echo "The setting must appear in $MENABITCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Bitcoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$MENABITCOIND_PIDDIR"
    chmod 0755 "$MENABITCOIND_PIDDIR"
    chown $MENABITCOIND_USER:$MENABITCOIND_GROUP "$MENABITCOIND_PIDDIR"
    chown $MENABITCOIND_USER:$MENABITCOIND_GROUP "$MENABITCOIND_CONFIGFILE"
    chmod 0660 "$MENABITCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$MENABITCOIND_PIDFILE" \
    --chuid $MENABITCOIND_USER:$MENABITCOIND_GROUP \
    --exec "$MENABITCOIND_BIN" \
    -- \
    -pid="$MENABITCOIND_PIDFILE" \
    -conf="$MENABITCOIND_CONFIGFILE" \
    -datadir="$MENABITCOIND_DATADIR" \
    -disablewallet \
    -daemon

